# 4g-5g-core-configuration/templates/cp/systemd/4G/zcore-ogstun.service
[Unit]
Description=Zcore network config
After=networking.service

[Service]
Type=simple

ExecStart=/opt/script/ogstun.sh

[Install]
WantedBy=multi-user.target
